# systemd-service/emqtt-check-cluster.timer
[Unit]
Description=Check EMQTT Cluster status each 15 sec

[Timer]
OnBootSec=2min
OnUnitActiveSec=15sec
Unit=emqtt-check-cluster.service
